# deploy/lm-studio.service.d/load-models.conf
[Service]
# Загрузка Ministral + nomic после старта LM Studio API
ExecStartPost=/bin/bash -c 'sleep 12 && /var/www/vladimirov/deploy/load-planrun-models.sh'
