# deploy/lm-studio.service.d/load-models.conf
[Service]
# lms server start завершается после запуска llmster — держим сервис активным через sleep infinity
ExecStart=
ExecStart=/bin/bash -c '/home/st_benni/.lmstudio/bin/lms server start; exec sleep infinity'
# Загрузка только embedding-модели после старта LM Studio API
ExecStartPost=/bin/bash -c 'sleep 12 && /var/www/planrun/deploy/load-planrun-models.sh'
