# deploy/planrun-ai.service.d/embed-model-fix.conf
[Service]
# Исправить ID модели эмбеддингов (без :2) — совпадает с LM Studio
Environment=LMSTUDIO_EMBED_MODEL=text-embedding-nomic-embed-text-v1.5
